# templates/supervisord.conf.j2
; Supervisor config file. Generated by Ansible role during deployment
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

{% if supervisord_unix_http_server %}
[unix_http_server]
{% for key, value in supervisord_unix_http_server.items() %}
{{ key }}={{ value }}
{% endfor %}
{% endif %}

{% if supervisord_inet_http_server %}
[inet_http_server]
{% for key, value in supervisord_inet_http_server.items() %}
{{ key }}={{ value }}
{% endfor %}
{% endif %}

{% if supervisord_supervisord %}
[supervisord]
{% for key, value in supervisord_supervisord.items() %}
{{ key }}={{ value }}
{% endfor %}
{% endif %}

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

{% if supervisord_supervisorctl %}
[supervisorctl]
{% for key, value in supervisord_supervisorctl.items() %}
{{ key }}={{ value }}
{% endfor %}
{% endif %}


; Programs
; ---------------------------------------------------------------
{% for program in supervisord_programs %}
[program:{{ program.name }}]
{% for key, value in program.items() %}
{% if key != 'name' %}
{{ key }}={{ value }}
{% endif %}
{% endfor %}
{% endfor %}


; Event listeners
; ---------------------------------------------------------------
{% for event in supervisord_eventlisteners %}
[eventlistener:{{ event.name }}]
{% for key, value in event.items() %}
{% if key != 'name' %}
{{ key }}={{ value }}
{% endif %}
{% endfor %}
{% endfor %}


; Groups
; ---------------------------------------------------------------

{% for group in supervisord_groups %}
[group:{{ group.name }}]
{% for key, value in group.items() %}
{% if key != 'name' %}
{{ key }}={{ value }}
{% endif %}
{% endfor %}
{% endfor %}

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

;[include]
;files = relative/directory/*.ini
